# tutorials/devops_for_small_to_medium_web_applications/sample-app/version4/infrastructure/15_certman/05_image/resources/certificate-updater.service
[Unit]
Description=certificate-updater
After=syslog.target
After=network.target
After=ossfs.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/python2.7 /opt/certificate-updater.py
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=certificate-updater

[Install]
WantedBy=multi-user.target